# laradock/workspace/crontab/laradock
* * * * * laradock php /var/www/citylit schedule:run >> /dev/null 2>&1
