# laradock/php-worker/supervisord.conf
[supervisord]
nodaemon=true
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=/var/www/run-queue.sh
autostart=true
autorestart=true
numprocs=8
redirect_stderr=true
