# laradock/deployments/citylit/worker.conf
[program:worker]
process_name=%(program_name)s_%(process_num)02d
command=php /home/vagrant/xmovement/artisan queue:work database --queue default,emails --sleep=3 --tries=3
autostart=true
autorestart=true
user=root
numprocs=8
redirect_stderr=true
